# share/ConfigShellUpgradeBySystemd/configshell-upgrade.service
[Unit]
Description="ConfigShell Upgrade Reposiroty"
Requires=configshell-upgrade.timer

[Service]
Type=simple
ExecStart=/
